# ignite/etc/sv/tap0/run
#!/bin/sh

sv once .

exec vde_switch -tap ${PWD##*/} -daemon
