# ignite/etc/sv/netctl-home/run
#!/bin/sh

/usr/lib/network/network start ${PWD##*netctl-}
chpst -b ${PWD##*/} pause
