# ignite/etc/sv/freenet/run
#!/bin/sh

sv once .
exec chpst -u freenet /opt/freenet/run.sh start
