# ignite/etc/sv/gogoc/run
#!/bin/sh

exec /opt/gogoc/bin/gogoc -n -y -f /opt/gogoc/bin/gogoc.conf 
