# ignite/etc/sv/oidentd/run
#!/bin/sh

exec /usr/bin/oidentd -i -u nobody -g nobody
