# ignite/etc/sv/colord/run
#!/bin/sh
exec /usr/lib/colord/colord 2>&1
