# ignite/etc/sv/colord-sane/run
#!/bin/sh
exec /usr/lib/colord/colord-sane 2>&1
