# ignite/etc/sv/znc/run
#!/bin/sh

exec chpst -u znc /usr/bin/znc -f
