# ignite/etc/sv/libvirtd/run
#!/bin/sh
. /etc/conf.d/libvirtd
exec /usr/bin/libvirtd $LIBVIRTD_ARGS
