# ignite/etc/sv/minbif/run
#!/bin/sh

exec chpst -u minbif /usr/bin/minbif /etc/minbif/minbif.conf
