# ignite/etc/sv/radvd/run
#!/bin/sh
exec /usr/sbin/radvd --nodaemon --logmethod=syslog
