# ignite/etc/sv/unbound/run
#!/bin/sh

exec /usr/sbin/unbound -d 2>&1
