# ignite/etc/sv/lxc-test-container/finish
#!/bin/sh

exec lxc-stop -n ${PWD##*lxc-}
