# ignite/etc/sv/lxc-test-container/run
#!/bin/sh

sv once .

exec lxc-start -d -n ${PWD##*lxc-}
