# ignite/etc/sv/netctl-xenbridge-dhcp/finish
#!/bin/sh

/usr/lib/network/network stop ${PWD##*netctl-}
