# ignite/etc/sv/netctl-xenbridge-dhcp/run
#!/bin/sh

sv once .
/usr/lib/network/network start ${PWD##*netctl-}
