# ignite/etc/sv/fms/run
#!/bin/sh

cd /opt/fms
exec chpst -u freenet -n 10 /opt/fms/fms
