# ignite/etc/sv/xinetd/run
#!/bin/sh

exec /usr/bin/xinetd -dontfork
