# ignite/etc/sv/qingy-generic/run
#!/bin/sh
tty=${PWD##*-}
exec env TERM=linux /usr/sbin/qingy $tty -t -n
